# systemd/user/qwen-nvfp4-vllm.service.in
[Unit]
Description=Qwen3.6-35B-A3B-NVFP4 vLLM OpenAI-compatible server
After=network-online.target docker.service
Conflicts=ornith-vllm.service mistral-medium-vllm.service

[Service]
Type=simple
WorkingDirectory={{INSTALL_ROOT}}
Environment=SPARK_MODEL_DIR={{MODEL_DIR}}
Environment=SPARK_VLLM_IMAGE={{VLLM_IMAGE}}
ExecStart={{INSTALL_ROOT}}/bin/qwen-nvfp4-vllm.sh
ExecStop=-{{INSTALL_ROOT}}/bin/docker-command.sh stop qwen-nvfp4-vllm
Restart=on-failure
RestartSec=15
TimeoutStartSec=900
TimeoutStopSec=60

[Install]
WantedBy=default.target
